# systemd/opsdoctor-agent.service
[Unit]
Description=OpsDoctor one-shot diagnostics collection
Documentation=https://github.com/opsdoctor/opsdoctor
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/opsdoctor-agent run
